# systemd/server-audit.service
[Unit]
Description=Server Audit Daily Report
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/server-audit/audit.py full
User=root
Nice=10
IOSchedulingClass=best-effort
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
